# salt/remotefirst/files/remotefirst_supervisor.ini
[program:remotefirst]
command = {{ salt['pillar.get']('remotefirst:site:venv') }}/bin/gunicorn remotefirst:app -w 4
directory = {{ salt['pillar.get']('remotefirst:site:root') }}
user = {{ salt['pillar.get']('remotefirst:site:user') }}
stdout_logfile = /var/log/gunicorn/remotefirst_gunicorn_stdout.log
stderr_logfile = /var/log/gunicorn/remotefirst_gunicorn_stderr.log
redirect_stderr = True
environment = PRODUCTION=1
